# backend/systemd.service
[Unit]
Description=Browser Agent Flask API
After=network.target

[Service]
Type=notify
User=www-data
Group=www-data
WorkingDirectory=/path/to/backend
Environment="PATH=/path/to/venv/bin"
Environment="FLASK_ENV=production"
Environment="SERVER_TYPE=gunicorn"
Environment="GEMINI_API_KEY=your_api_key_here"
Environment="SECRET_KEY=your_secret_key_here"
ExecStart=/path/to/venv/bin/gunicorn --config gunicorn_config.py "flask_app.app:create_app()"
ExecReload=/bin/kill -s HUP $MAINPID
KillMode=mixed
TimeoutStopSec=5
PrivateTmp=true
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target

